# files/etc/systemd/system/backup.service
[Service]
Type=oneshot
Restart=on-failure

ExecStart=/usr/local/bin/backup.sh
