# files/etc/systemd/system/scrub.service
[Service]
Type=oneshot
Restart=on-failure

ExecStart=/usr/bin/btrfs scrub start -B /
