# files/etc/systemd/system/sync.service
[Service]
Type=oneshot

User=kai
ExecStart=/usr/local/bin/sync.sh
