# cluster/applications/chrome-devtools-protocol-server/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0

[program:xvfb]
command=Xvfb :0 -screen 0 1920x1080x24
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

[program:chromium]
command=chromium --remote-debugging-port=9222 --no-sandbox --disable-gpu
environment=DISPLAY=":0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# https://issues.chromium.org/issues/41487252
[program:socat]
command=socat TCP-LISTEN:59222,fork,reuseaddr TCP:127.0.0.1:9222
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

[unix_http_server]
file=/tmp/supervisor.sock
