# files/etc/systemd/user/claude-daily-task.service
[Unit]

[Service]
Type=oneshot

ExecStart=/usr/local/bin/claude-daily-task.sh
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal
