# files/etc/profile.d/android.sh
export ANDROID_HOME="/home/${USER}/Android/Sdk"
PATH="${ANDROID_HOME}/platform-tools:${ANDROID_HOME}/emulator:$PATH"
