# files/etc/profile.d/gemini.sh
export GEMINI_SANDBOX_IMAGE=ghcr.io/hippocampus-dev/hippocampus/gemini-sandbox:main
