# files/etc/profile.d/path.sh
export PATH="/home/${USER}/bin:/home/${USER}/.local/bin:/home/${USER}/go/bin:/home/${USER}/.cargo/bin:/home/${USER}/llm:/opt/intellij-idea-ultimate-edition/bin:${PATH}"
