# files/etc/profile.d/path.sh
export PATH="/home/kai/bin:/home/kai/.local/bin:/home/kai/.cargo/bin:/home/kai/llm:/opt/intellij-idea-ultimate-edition/bin:/opt/google-cloud-cli/bin:$PATH"
