# files/etc/profile.d/sslkeylog.sh
export SSLKEYLOGFILE="${XDG_RUNTIME_DIR:-/tmp}/sslkey.log"
